# deb_build/DEBIAN/prerm
#!/bin/bash
set -e

echo "[SuperManager] Zastavuji službu…"
systemctl stop supermanager.service || true
systemctl disable supermanager.service || true
